# launchd/com.art.enforce_exchange_alerts.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.art.enforce_exchange_alerts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/REPLACE_WITH_YOUR_USERNAME/bin/enforce_exchange_alerts</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/Library/Logs/enforce_exchange_alerts.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/Library/Logs/enforce_exchange_alerts.log</string>
</dict>
</plist>
